# speed_collector_autostart.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2025-10-22T00:00:00</Date>
    <Author>vd-speed-test</Author>
    <Description>Run Internet Speed Collector every 15 minutes (silent)</Description>
  </RegistrationInfo>

  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2025-10-22T00:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
      <Repetition>
        <Interval>PT15M</Interval>
        <Duration>P1D</Duration>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
    </CalendarTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <!-- Replace with your Windows user. Examples:
           varad
           DESKTOP-ABCD123\varad
           varad@yourdomain.com
      -->
      <UserId>varad-studio\varad</UserId>

      <!-- Password = "Run whether user is logged on or not"
           (Windows will ask for password when you register/import) -->
      <LogonType>S4U</LogonType>

      <!-- Avoid UAC/UI issues; only use HighestAvailable if truly needed -->
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <MultipleInstancesPolicy>StopExisting</MultipleInstancesPolicy>
    <ExecutionTimeLimit>PT10M</ExecutionTimeLimit>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>false</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>

    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>

    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>true</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <Priority>4</Priority>
  </Settings>

  <Actions Context="Author">
    <Exec>
      <Command>C:\Windows\System32\wscript.exe</Command>
      <Arguments>"C:\vd-speed-test\run_speed_collector.vbs"</Arguments>
      <WorkingDirectory>C:\vd-speed-test</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
